# scripts/myspeed.rc
#!/sbin/openrc-run

# OpenRC metadata
description="myspeed service"
command="/usr/local/bin/node"
command_args="server"
command_user="myspeed:myspeed"
directory="/opt/myspeed"
name="myspeed"

# Logging (optional)
pidfile="/run/myspeed.pid"
output_log="/var/log/myspeed.log"
error_log="/var/log/myspeed.err"

depend() {
    need net
}

start_pre() {
    ebegin "Preparing to start myspeed"

    # Ensure /opt/myspeed exists
    if [ ! -d "$directory" ]; then
        eerror "The $directory directory does not exist!"
        return 1
    fi

    # Ensure bin/ and data/ directories exist inside /opt/myspeed
    for dir in bin data; do
        if [ ! -d "$directory/$dir" ]; then
            einfo "Creating $directory/$dir directory"
            mkdir -p "$directory/$dir" || return 1
        fi
    done

    # Set ownership of all files/folders to myspeed:myspeed
    einfo "Setting ownership of $directory to myspeed:myspeed"
    chown -R myspeed:myspeed "$directory" || return 1

    eend 0
}

stop() {
    ebegin "Stopping myspeed"
    start-stop-daemon --stop --quiet --pidfile "$pidfile"
    eend $?
}

start() {
    ebegin "Starting myspeed"
    start-stop-daemon --start --quiet --pidfile "$pidfile" \
                      --make-pidfile --chdir "$directory" \
                      --user "$command_user" \
                      --env NODE_ENV=production \
                      --exec "$command" -- $command_args >> "$output_log" 2>> "$error_log" &
    eend $?
}
