# scripts/myspeed.rc
#!/sbin/openrc-run

# OpenRC metadata
name="myspeed"
description="myspeed service"
command="/usr/bin/node"
directory="/opt/myspeed"
command_args=$directory/server
command_user="myspeed:myspeed"
command_background=true

# Logging (optional)
pidfile="/run/myspeed.pid"
#output_log="/var/log/myspeed.log"
#error_log="/var/log/myspeed.err"

depend() {
    need net
}

start_pre() {
    ebegin "Preparing to start myspeed"

    # Ensure /opt/myspeed exists
    if [ ! -d "$directory" ]; then
        eerror "The $directory directory does not exist!"
        return 1
    fi

    # Ensure bin/ and data/ directories exist inside /opt/myspeed
    for dir in bin data; do
        if [ ! -d "$directory/$dir" ]; then
            einfo "Creating $directory/$dir directory"
            mkdir -p "$directory/$dir" || return 1
        fi
    done

    # Set ownership of all files/folders to myspeed:myspeed
    einfo "Setting ownership of $directory to $command_user"
    chown -R "$command_user" "$directory" || return 1

	export NODE_END=production

    eend 0
}
